# launchd/com.empp.agent.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  launchd template for the book-update build-agent runner.

  Install:
    1. Edit the paths and env vars below (NODE path, WorkingDirectory, secrets).
    2. cp launchd/com.empp.agent.plist ~/Library/LaunchAgents/com.empp.agent.plist
    3. launchctl load -w ~/Library/LaunchAgents/com.empp.agent.plist
    4. Logs: tail -f /tmp/empp-agent.out /tmp/empp-agent.err

  Stop:    launchctl unload -w ~/Library/LaunchAgents/com.empp.agent.plist

  Do NOT commit real secrets — fill AGENT_ENROLL_SECRET / APPLE_* on the agent host only.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.empp.agent</string>

  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/CHANGEME/01dev/electron-multi-platform-packager/src/agent/runner.js</string>
  </array>

  <key>WorkingDirectory</key>
  <string>/Users/CHANGEME/01dev/electron-multi-platform-packager</string>

  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin:/Applications/Xcode.app/Contents/Developer/usr/bin</string>

    <key>BOOKUPDATE_API</key>
    <string>https://akillitahta.ndr.ist/api/v1</string>

    <key>PACKAGER_API</key>
    <string>http://127.0.0.1:3001</string>

    <key>AGENT_CAPS</key>
    <string>android,macos</string>

    <!-- Required for first-time enrollment only (then a token file is reused). -->
    <key>AGENT_ENROLL_SECRET</key>
    <string>CHANGEME</string>

    <!-- macOS signing (best-effort; omit to ship unsigned). -->
    <key>APPLE_SIGN_IDENTITY</key>
    <string>Developer ID Application: Mehmet Arslan (335PPR74QM)</string>
    <key>APPLE_TEAM_ID</key>
    <string>335PPR74QM</string>
    <!-- Prefer a stored notarytool keychain profile: -->
    <key>APPLE_NOTARY_PROFILE</key>
    <string>empp-notary</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>

  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>/tmp/empp-agent.out</string>
  <key>StandardErrorPath</key>
  <string>/tmp/empp-agent.err</string>
</dict>
</plist>
